# nmpcron.cron
*/5 * * * * /usr/local/bin/npm run archiveAppointments >> /tmp/txt.log
