# nmpcron.cron
*/5 * * * * bash -c 'npm run archiveAppointments' >> /tmp/cron.log 2>&1
