# nmpcron.cron
*/5 * * * * node /usr/local/bin/npm run archiveAppointments >> /var/log/cron.log 2>&1
