# nmpcron.cron
*/5 * * * * npm run archiveAppointments
